# deploy/vm/token-trackr-worker.service
[Unit]
Description=Token Trackr Worker Service (Scheduler)
Documentation=https://github.com/your-org/token-trackr
After=network.target postgresql.service token-trackr-api.service
Wants=postgresql.service

[Service]
Type=exec
User=tokentrackr
Group=tokentrackr
WorkingDirectory=/opt/token-trackr

# Environment
EnvironmentFile=/etc/token-trackr/token-trackr.env
Environment=SCHEDULER_ENABLED=true

# Start command
ExecStart=/opt/token-trackr/venv/bin/python -m backend.jobs.scheduler

# Restart configuration
Restart=always
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=token-trackr-worker

# Resource limits
LimitNOFILE=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Allowed paths
ReadWritePaths=/var/log/token-trackr /opt/token-trackr/reports

[Install]
WantedBy=multi-user.target

